# Library/dpkg/info/openssh-server.list
/.
/var
/var/jb
/var/jb/Library
/var/jb/Library/LaunchDaemons
/var/jb/Library/LaunchDaemons/com.openssh.sshd.plist
/var/jb/etc
/var/jb/etc/pam.d
/var/jb/etc/pam.d/sshd
/var/jb/etc/ssh
/var/jb/etc/ssh/moduli
/var/jb/etc/ssh/sshd_config
/var/jb/usr
/var/jb/usr/libexec
/var/jb/usr/libexec/sshd-keygen-wrapper
/var/jb/usr/sbin
/var/jb/usr/sbin/sshd
/var/jb/usr/share
/var/jb/usr/share/doc
/var/jb/usr/share/doc/openssh-server
/var/jb/usr/share/doc/openssh-server/README.zst
/var/jb/usr/share/doc/openssh-server/TODO.zst
/var/jb/usr/share/man
/var/jb/usr/share/man/man5
/var/jb/usr/share/man/man5/moduli.5.zst
/var/jb/usr/share/man/man5/sshd_config.5.zst
/var/jb/usr/share/man/man8
/var/jb/usr/share/man/man8/sshd.8.zst
